# docker/mt5-terminal/supervisord.conf
; supervisord configuration for the mt5-terminal container.
;
; This file controls the four long-running processes inside the
; Wine + MT5 + RPyC-bridge container:
;   1. Xvfb            (headless X server at :99)
;   2. x11vnc          (VNC server bound to 127.0.0.1:5900 by compose)
;   3. websockify+novnc (browser-accessible VNC on 6080)
;   4. mt5terminal.exe (background Windows process under Wine)
;   5. python mt5_bridge_server.py (RPyC server, port 18812)
;
; Process ordering matters: Xvfb must be up before any wine process.
; supervisord's priority field controls start order (lower = earlier).
; We use 10/20/30/40/50 so restarts are deterministic.

[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
environment=DISPLAY=":99",WINEPREFIX="/root/.wine",WINEARCH="win64",WINEDEBUG="-all"

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ============================================================ Xvfb
;
; Headless X server on :99. Wine's mt5terminal.exe is a Windows GUI
; process — it needs *some* X server, even if we never display it.
; xvfb is the canonical choice.
[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1280x720x24
priority=10
autostart=true
autorestart=true
startsecs=2
startretries=3
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.err

; ============================================================ x11vnc
;
; VNC server bound to 127.0.0.1 (security: see AGENTS.md §4h-5).
; The actual bind address is enforced by docker-compose.prod.yml's
; `127.0.0.1:5900:5900` port mapping, but we also pin x11vnc to
; localhost as defense in depth.
[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -shared -rfbport 5900
    -localhost -nopw -noipv6
    -xkb
priority=20
autostart=true
autorestart=true
startsecs=2
startretries=3
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.err

; ============================================================ novnc
;
; WebSocket-to-VNC bridge. Users can hit http://127.0.0.1:6080/vnc.html
; in a browser to get a full VNC session (one-time Vantage account
; login). Loopback only via compose.
[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc 6080 localhost:5900
priority=30
autostart=true
autorestart=true
startsecs=2
startretries=3
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.err

; ============================================================ mt5terminal
;
; MetaTrader 5 terminal under Wine. Background only — it stays open
; to keep the active session alive. Note: we do NOT pass /portable
; because the persistent volume mount in compose holds the terminal
; state across container restarts. If you ever need a clean
; re-install, delete the mt5-data volume.
[program:mt5terminal]
command=/usr/bin/xvfb-run -a wine /root/.wine/drive_c/Program Files/MetaTrader 5/terminal64.exe /portable
priority=40
autostart=true
autorestart=true
startsecs=10
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/mt5terminal.log
stderr_logfile=/var/log/supervisor/mt5terminal.err

; ============================================================ mt5_bridge
;
; The RPyC server. The actual MT5 connector. It must start AFTER
; mt5terminal.exe so that mt5.initialize() finds the terminal.
; We sleep 15s before launch to give mt5terminal time to draw its
; main window (which is when the terminal becomes connectable).
[program:mt5_bridge]
command=/bin/sh -c "sleep 15 && exec /usr/bin/xvfb-run -a wine python /opt/mt5-bridge/mt5_bridge_server.py --port %(ENV_MT5_BRIDGE_PORT)s --host %(ENV_MT5_BRIDGE_HOST)s"
priority=50
autostart=true
autorestart=true
startsecs=30
startretries=10
stopasgroup=true
killasgroup=true
environment=PYTHONPATH="Z:\\opt\\mt5-bridge"
stdout_logfile=/var/log/supervisor/mt5_bridge.log
stderr_logfile=/var/log/supervisor/mt5_bridge.err
